# service/dock.service
[Unit]
Description="LORIS Dock"

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
ExecStart=/home/root/launch-dock.sh
