# service/ipc.service
[Unit]
Description="LORIS IPC daemon"

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
ExecStart=/home/root/launch-dock.sh
